# deploy/robot-nl.service
[Unit]
Description=Robot Natural Language Controller
Documentation=https://github.com/your-org/robot_nl_controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=robot
Group=robot
WorkingDirectory=/opt/robot_nl_controller

# GCP 資格情報・設定（secrets.env に記載）
EnvironmentFile=/etc/robot_nl/secrets.env

# ROS2 環境のセットアップ
ExecStartPre=/bin/bash -c "source /opt/ros/humble/setup.bash"
ExecStart=/bin/bash -c "source /opt/ros/humble/setup.bash && python3 main.py /etc/robot_nl/config.yaml"

# クラッシュ時に自動再起動
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# 標準出力をジャーナルに記録
StandardOutput=journal
StandardError=journal
SyslogIdentifier=robot-nl

[Install]
WantedBy=multi-user.target
